# supervisor/eq-horizon.ini
[program:eq-horizon]
process_name=%(program_name)s
command=php {ARTISAN_PATH} horizon
autostart=true
autorestart=true
user=apache
redirect_stderr=true
stdout_logfile=/data/eq-data/logs/horizon.log
stopwaitsecs=28860